# services/applications/mitk-workbench/docker/files/supervisord-mitk.conf
[program:mitk]
priority=50
process_name=mitk
command=x-terminal-emulator -e "/mitk/startMITK.sh | tee $HOME/mitk-logfile"
user=%USER%
environment=DISPLAY=":1",HOME="%HOME%",USER="%USER%"