# services/applications/slicer-workbench/docker/files/supervisor-slicer.conf
[program:slicer]
priority=50
process_name=Slicer
command=/opt/slicer/Slicer
user=%USER%
environment=DISPLAY=":1",HOME="%HOME%",USER="%USER%"
