# data-processing/processing-pipelines/mitk-flow/processing-containers/docker/files/supervisord-mitk-flowbench.conf
[program:mitk-flow]
priority=50
process_name=mitk-flow
command=x-terminal-emulator -e "/mitk-flow/startMITKFlow.sh | tee $HOME/logfile"
user=%USER%
environment=DISPLAY=":1",HOME="%HOME%",USER="%USER%"