# data-processing/processing-pipelines/m2olie/processing-containers/m2olie-workbench/docker/files/supervisord-m2olie.conf
[program:mitk]
priority=50
process_name=mitk
command=x-terminal-emulator -e "/m2olie/startMITK.sh | tee $USER/logfile"
user=%USER%
environment=DISPLAY=":1",HOME="%HOME%",USER="%USER%"